# usb-umount@.service
[Unit]
Description=Unmount USB drive %i
After=systemd-udevd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb-umount.sh /dev/%i
